# templates/cdnjs-bot2.timer
[Unit]
Description=Run cdnjs-bot2.service

[Timer]
Unit=cdnjs-bot2.service
# every 3 mins
OnUnitActiveSec=3min

[Install]
WantedBy=multi-user.target
