# templates/cdnjs-bot.timer
[Unit]
Description=Run cdnjs-bot.service

[Timer]
Unit=cdnjs-bot.service
# every 60 mins
OnUnitActiveSec=60min

[Install]
WantedBy=multi-user.target
