# templates/update-website.timer
[Unit]
Description=Run update-website.service

[Timer]
# every day at 7am
OnCalendar=*-*-* 7:00:00

[Install]
WantedBy=timers.target
